# systemd/actions-runner.service.d/resilience.conf
# OOM protection and auto-restart for the GitHub Actions runner.
#
# Deployed to:
#   /etc/systemd/system/actions.runner.hertie-data-science-lab.ds01-runner.service.d/resilience.conf
#
# Context: The runner was OOM-killed on 2026-03-27 because user workloads
# exhausted system memory and the .NET runtime's large virtual address space
# inflated its oom_score. Without Restart= it stayed dead for 4 days.

[Service]
Restart=always
RestartSec=10
OOMScoreAdjust=-500
