# _legacy/deploy/ds01-api.service
[Unit]
Description=DS01 Job Submission API
Documentation=https://github.com/YOUR_ORG/ds01-infra
After=network.target
Wants=network.target

[Service]
Type=exec
User=root
WorkingDirectory=/opt/ds01-infra/scripts/api
ExecStartPre=/bin/mkdir -p /var/lib/ds01/api
# Bind address: 127.0.0.1:8765 — localhost only, Cloudflare Tunnel proxies inbound
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 127.0.0.1 --port 8765 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ds01-api

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/ds01/api /var/log/ds01
ProtectHome=yes

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
