# systemd/ds01-runner.service
[Unit]
Description=DS01 Job Runner
Documentation=file:///opt/ds01-jobs/README.md
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=ds01
Group=ds01
WorkingDirectory=/opt/ds01-jobs
# Venv at /var/lib/ds01-jobs/venv is created/refreshed by ExecStartPre as ds01.
# Uses system uv at /usr/local/bin/uv and uv-managed Python 3.13 at
# /usr/local/share/uv-python/ (both world-readable, no dependency on user homes).
ExecStartPre=/usr/local/bin/uv sync --locked --all-groups --project /opt/ds01-jobs
ExecStart=/var/lib/ds01-jobs/venv/bin/ds01-job-runner
Environment=UV_PROJECT_ENVIRONMENT=/var/lib/ds01-jobs/venv
Environment=UV_PYTHON_INSTALL_DIR=/usr/local/share/uv-python
Environment=UV_PYTHON=/usr/local/share/uv-python/cpython-3.13-linux-x86_64-gnu/bin/python3.13
Environment=UV_CACHE_DIR=/var/lib/ds01-jobs/uv-cache
Restart=always
RestartSec=5
UMask=0002
KillMode=process
EnvironmentFile=/etc/ds01-jobs/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ds01-runner

# Hardening (relaxed - needs Docker/sudo access)
ProtectSystem=full
ReadWritePaths=/opt/ds01-jobs/data /var/log/ds01 /var/lib/ds01-jobs
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
